# lovestory.service
[Unit]
Description=LoveStory Node.js Backend Service
Documentation=https://github.com/nie11kun/LoveStory
After=network.target

[Service]
Type=simple
# 请根据实际运行本应用的用户修改，建议不要使用 root
User=root

# 必须修改为你服务器上 LoveStory 源码的完全绝对物理路径
WorkingDirectory=/path/to/your/LoveStory

# 如果你在使用 nvm 或者其他 node 版本管理器，请务必使用 node 的绝对路径（通过 which node 查询）
ExecStart=/usr/bin/node server.js

Restart=on-failure
RestartSec=10

# 生产环境标记
Environment=NODE_ENV=production

# 標準輸出與錯誤輸出 (可選)
# StandardOutput=syslog
# StandardError=syslog
# SyslogIdentifier=lovestory

[Install]
WantedBy=multi-user.target
